# mount/mount.nfs
#!/bin/bash
### BEGIN INIT INFO 
# @provides		mount.store
# @author		Antonio Membrides Espinosa
# @author_URI:  https://github.com/ameksike
# @package  	Task
# @date			15/01/2018
# @update		15/02/2018
# @version    	1.0
# @license    	GPL
# @short		
# @description 	Mount FTP Server
### END INIT INFO

mount -t nfs 192.168.3.233:/media/raid0 /media/store
